# startup/dr-provision.sysv
#!/bin/bash
# chkconfig: 2345 90 10
# description: DigitalRebar Provision Provisioner and DHCP server
PATH="$PATH:/usr/local/bin"
touch /var/log/dr-provision.log
case $1 in
    start)
        /usr/local/bin/dr-provision &>>/var/log/dr-provision.log </dev/zero &
        disown;;
    stop)
        pkill dr-provision;;
    status)
        pgrep dr-provision;;
esac

        
